# tomcat-mongo/conf/mongo.conf
[program:mongod]
command=/usr/bin/mongod --port 27017 --quiet 
process_name=%(program_name)s
autostart=true
user=root
startsecs=3
redirect_stderr=true
stdout_logfile=/var/log/mongod.log
